# deploy/depth-recog.service
[Unit]
Description=Rider depth estimation — DepthAnything-V2-Small (pour bullet-time)
After=network.target

[Service]
Type=simple
User=ben
Group=ben
WorkingDirectory=/opt/rider-recognition

Environment=STREAM_URL=http://localhost:8810/stream.mjpeg
Environment=OUTPUT_NPY=/tmp/rider-depth.npy
Environment=OUTPUT_JSON=/tmp/rider-depth.json
Environment=DEPTH_PERIOD=0.5
Environment=GPU_ID=0
EnvironmentFile=-/etc/rider-recognition/depth-recog.env

ExecStart=/opt/rider-recognition/venv/bin/python3 /opt/rider-recognition/scripts/depth_recog_service.py

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
